# systemd/.config/systemd/user/battery-level.service
[Unit]
Description=Battery level notification checker
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/scripts/battery-level.sh
